# config/autostart/middle-click-hold.desktop
[Desktop Entry]
Type=Application
Name=Middle Click Hold
Comment=Hold middle mouse button to hold spacebar (voice input)
Exec=python3 $HOME/.local/bin/middle-click-hold.py
X-GNOME-Autostart-enabled=true
Hidden=false
NoDisplay=true
