# config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=Xbindkeys
Exec=xbindkeys
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
